# goddadyDynDNS.timer
[Unit]
Description=Dynamic DNS for godaddy timer file

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
